# qt_app/systemd/fitpro-ultima-kiosk.service
[Unit]
Description=FitPro Ultima Laser native Qt kiosk
After=graphical.target hairkiller-backend.service
Wants=hairkiller-backend.service

[Service]
Type=simple
User=jetson
WorkingDirectory=/opt/fitpro-ultima-laser/qt_app
Environment=DISPLAY=:0
Environment=QT_QPA_PLATFORM=xcb
Environment=FITPRO_API_BASE_URL=http://127.0.0.1:8000
ExecStart=/opt/fitpro-ultima-laser/qt_app/run-kiosk.sh
Restart=always
RestartSec=2

[Install]
WantedBy=graphical.target
